# roles/web/templates/supervisor.atlas.conf
[program:uwsgi]
command=uwsgi /etc/uwsgi/apps-enabled/atlas.ini
user = {{atlas_webuser}}
group = {{atlas_webuser}}
autostart={{atlas_environment == 'prod'}}
startsecs=2
stopsignal=INT
